# bin/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
#::respawn:-/bin/sh
#::respawn:/sbin/getty -L ttySAC0 115200 vt100

#::respawn:-/bin/login
#::restart:/sbin/init

# Start an "askfirst" shell on the console (whatever that may be)
::askfirst:-/bin/login

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r 2>/dev/null

